# debian/postinst
#!/bin/sh

update-rc.d wb-mqtt-apcsnmp defaults
invoke-rc.d wb-mqtt-apcsnmp start
